# contrib/systemd/mako.service.in
[Unit]
Description=A lightweight Wayland notification daemon
Documentation=man:mako(1)
PartOf=wayland-session.target

[Service]
Type=dbus
BusName=org.freedesktop.Notifications
ExecStart=@bindir@/mako

[Install]
WantedBy=wayland-session.target
