# mako.systemd.service.in
[Unit]
Description=Lightweight Wayland notification daemon
Documentation=man:mako(1)
PartOf=graphical-session.target

[Service]
Type=dbus
BusName=org.freedesktop.Notifications
ExecCondition=/bin/sh -c '[[ "$XDG_SESSION_TYPE" == "wayland" ]]'
ExecStart=@bindir@/mako
ExecReload=@bindir@/mako reload

[Install]
WantedBy=graphical-session.target
